# falcon/supervisor.conf
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisor/supervisord.log;
;stdout_logfile=/dev/stdout
;stdout_logfile_maxbytes=0
;stderr_logfile=/dev/stderr
;stderr_logfile_maxbytes=0
loglevel = debug

[program:app]
redirect_stderr = true
autorestart = true
;command = gunicorn  --chdir /app/api -w 2  -b 0.0.0.0:80 main:api
command = uwsgi --wsgi-file /app/app.py --gevent 100 --socket :3031 --callable app --http :9000
